# debian/postinst
#!/bin/sh
# postinst script for thruk
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package


case "$1" in
    configure)
        # enable modules we need
        a2enmod alias
        a2enmod fcgid
        a2enmod auth_basic
        a2enmod rewrite

        # set permissions
        chown www-data: /var/lib/thruk
        chown www-data: /var/cache/thruk
        chown www-data: /var/log/thruk
        chown 755 /usr/share/thruk/fcgid_env.sh

        # reload new apache config
        if which invoke-rc.d >/dev/null 2>&1; then
            invoke-rc.d apache2 reload
        else
            /etc/init.d/apache2 reload
        fi

        echo "Thruk has been configured for http://$(hostname)/thruk/. User and password is 'thrukadmin'."
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
