# debian/prerm
#!/bin/sh
# prerm   script for thruk
#
# see: dh_installdeb(1)

set -e

case "$1" in
  remove|deconfigure)
    /usr/bin/thruk -a uninstallcron --local || echo "ERROR: removing cron failed, please remove www-data cron manually"
    invoke-rc.d thruk stop || echo
    rm -f /usr/share/thruk/root/thruk/plugins/*
    rmdir /usr/share/thruk/root/thruk/plugins
  ;;

  upgrade|failed-upgrade)
  ;;

  *)
    echo "prerm called with unknown argument \`$1'" >&2
    exit 0
  ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
